# indexer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Infrastructure ────────────────────────────────────────────────────────────

[program:redis]
command=redis-server --save "" --appendonly no --loglevel warning
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nats]
command=nats-server -js -c /app/nats.conf
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Application ───────────────────────────────────────────────────────────────
; Both start after a short delay so Redis + NATS are accepting connections first.

[program:indexer]
command=sh -c "sleep 3 && /app/indexer index --config /app/configs/b2b/sui.yaml --chains=sui_testnet --from-latest"
autostart=true
autorestart=true
priority=20
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:dispatcher]
command=sh -c "sleep 4 && /app/dispatcher"
autostart=true
autorestart=true
priority=20
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
